# app-admin/supervisor/files/supervisord.initd
#!/sbin/runscript

description="Client/server system for controlling processes on UNIX-like operating systems"

command="/usr/bin/supervisord"
pidfile="/var/run/supervisord.pid"
command_args="--nodaemon --pidfile ${pidfile} ${SUPERVISORD_OPTIONS}"
command_background="yes"

required_files="/etc/supervisord.conf"
